# NoVNC/Context/app/conf.d/fluxbox.conf
[program:fluxbox]
command=bash -c "sleep 3; exec fluxbox"
autorestart=true
priority=100
startsecs=5
