# NoVNC/Context/app/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 8090 localhost:5900
autorestart=true
priority=200